# 99-monitor-detected.rules
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card0" ENV{HOTPLUG}=="1", RUN+="/etc/udev/scripts/monitor-detected-to-configure.sh"
